# services.d/php7.0-fpm/run
#!/bin/sh
exec /usr/bin/php-fpm -R --nodaemonize --force-stderr --fpm-config /etc/php/php-fpm.conf --pid /var/run/php-fpm.pid